# data-sync/src/main/resources/service/Inf_db_Service_Sync.timer
[Unit]
Description=Run service every 10 minutes

[Timer]
OnCalendar=00/8:00
